# dotfiles/systemd/user/waybar-updates.service
[Unit]
Description=Check for pacman updates and notify waybar

[Service]
Type=oneshot
ExecStart=%h/.config/waybar/updates.sh
